# gunicorn.conf
[program:gunicorn]
directory=/home/ubuntu/favorite-things/favourite_things_api
command=/home/ubuntu/.local/share/virtualenvs/favorite-things-NiYzl7a6/bin/gunicorn --workers 3  --bind unix:/home/ubuntu/favorite-things/app.sock favourite_things_api.wsgi:applica$
autostart=true
autorestart=true
stderr_logfile=/var/log/gunicorn/gunicorn.err.log
stdout_logfile=/var/log/gunicorn/gunicorn.out.log
